# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[eventlistener:failfast]
command=/usr/local/bin/supervisor-failfast
events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
buffer_size=20
autostart=true
autorestart=true
startsecs=0
priority=1
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:warp-svc]
command=/usr/local/bin/warp-svc-wrapper
autostart=true
autorestart=false
startsecs=2
priority=10
stopsignal=TERM
stopwaitsecs=30
killasgroup=true
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/dev/null
stderr_logfile_maxbytes=0

[program:warp-bootstrap]
command=/usr/local/bin/warp-bootstrap
autostart=true
autorestart=false
exitcodes=0
startsecs=0
priority=20
stopsignal=TERM
stopwaitsecs=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:warp-socat]
command=/usr/local/bin/warp-socat
autostart=true
autorestart=false
startsecs=1
priority=30
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
