# supervisord.conf
[supervisord]
nodaemon=false
logfile=/dev/null
pidfile=/run/supervisord.pid

[program:warp-svc]
command=/bin/sh -lc 'exec "$(command -v warp-svc)"'
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
